# 10-screen.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/bash -c '/home/miro/skriptz/display_cfg.sh'"

